# deploy/homelab/systemd/user/cortex-homelab-deploy-listener.service
[Unit]
Description=Cortex homelab deploy listener (host Docker control)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/greyhill/Documents/Cortex
Environment=NODE_ENV=production
Environment=CORTEX_DEPLOY_LISTENER_PORT=9092
Environment=CORTEX_DEPLOY_LISTENER_HOST=0.0.0.0
EnvironmentFile=-/home/greyhill/Documents/Cortex/deploy/homelab/.env
ExecStart=/home/greyhill/Documents/Cortex/scripts/homelab-deploy-listener-start.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
