# deploy/homelab/systemd/user/cortex-homelab-deploy-watch.timer
[Unit]
Description=Poll for Cortex homelab deploy every 2 minutes (git + source hash)

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min
Persistent=true

[Install]
WantedBy=timers.target
